# deploy/systemd/nostr-relay.service
[Unit]
Description=Bılbıl Nostr Relay (rrainn/nostr-relay)
After=network.target

[Service]
Type=simple
User=nostr
Group=nostr
WorkingDirectory=/opt/nostr/src/nostr-relay
ExecStart=/usr/bin/node dist/src/WebSocket.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
IPAddressDeny=any
IPAddressAllow=127.0.0.0/8
IPAddressAllow=::1/128
ReadWritePaths=/opt/nostr/src/nostr-relay

[Install]
WantedBy=multi-user.target
